# ubuntu/supervisor_program.conf
[program:program_name]
command=python script_name.py
directory=/opt/starting/directory
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
